# scripts/webcface-server.service.in
[Unit]
Description=WebCFace Server

[Service]
Type=exec
Restart=always
ExecStart=@prefix@/@bindir@/webcface-server

[Install]
WantedBy=multi-user.target
