# cmake/webcface-server.service
[Unit]
Description=WebCFace Server

[Service]
Type=exec
Restart=always
ExecStart=webcface-server

[Install]
WantedBy=multi-user.target
